# restic/.config/systemd/user/restic-work.service
[Unit]
Description=Backup experiments

[Service]
ExecStart=/home/nasser/.restic